# debian/lingmo-installer.postinst
#! /bin/sh
set -e

. /usr/share/debconf/confmodule

if [ -e /etc/rc2.d/S20lingmo-installer ]; then
	update-rc.d -f lingmo-installer remove >/dev/null 2>&1
fi

#DEBHELPER#

exit 0
